# icaros-lxde-light_0.1-1/etc/xdg/autostart/icaros.desktop
[Desktop Entry]
Type=Application
Name=Icaros
Comment=Start Icaros Desktop
Exec=bash -c "/usr/bin/start-icaros"
OnlyShowIn=LXDE
