# icaros-lxde-light_0.1-1/etc/skel/.config/autostart/icaros.desktop
[Desktop Entry]
Type=Application
Name=Icaros
Comment=Start Icaros Desktop
Exec=bash -c "sudo /usr/lib/icaros/S/mount_overlay.sh && $HOME/IcarosDesktop/start-icaros.sh"
OnlyShowIn=LXDE
