# icaros-light_2.2.9-1/DEBIAN/postrm
groupdel aros

systemctl disable arosnet
